# bin/syncd.service
[Unit]
Description=Service Sycnd
After=network.target

[Service]
Type=forking
ExecStart=/bin/sh -c '/opt/syncDirectory/syncd.sh' start
ExecStop=/bin/sh -c '/opt/syncDirectory/syncd.sh' stop
TimeoutStartSec=0

[Install]
WantedBy=default.target
