# bin/syncd.service
[Unit]
Description=Service Sycnd
After=network.target

[Service]
Type=forking
ExecStart=/opt/syncd/syncd.sh start
ExecStop=/opt/syncd/syncd.sh stop
TimeoutStartSec=0

[Install]
WantedBy=default.target
